# debian/chef-server-webui/debian/chef-server-webui/DEBIAN/postrm
#!/bin/sh
# postrm script for chef-server-webui
#
# see: dh_installdeb(1)

set -e

TEMPLATE=/usr/share/chef/webui.rb
CONFIGFILE=/etc/chef/webui.rb

case "$1" in
    remove)
    ;;
    purge)
      rm -rf /var/log/chef/server-webui*
      rm -rf /var/log/chef/merb-webui*
      test -x /usr/bin/ucf && ucf --purge $TEMPLATEFILE $CONFIGFILE
      rm -rf $CONFIGFILE $CONFIGFILE.ucf-old $CONFIGFILE.ucf-new $CONFIGFILE.ucf-dist
    ;;
    upgrade|failed-upgrade|abort-install|abort-upgrade|disappear|remove)
    ;;

    *)
        echo "postrm called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

# Automatically added by dh_installinit
if [ "$1" = "purge" ] ; then
	update-rc.d chef-server-webui remove >/dev/null
fi
# End automatically added section
# Automatically added by dh_installdebconf
if [ "$1" = purge ] && [ -e /usr/share/debconf/confmodule ]; then
	. /usr/share/debconf/confmodule
	db_purge
fi
# End automatically added section


exit 0
